# src/etc/supervisor/conf.d/ftp.conf
[program:ftp]
command=/usr/sbin/pure-ftpd -c 50 -C 10 -l puredb:/etc/pure-ftpd/pureftpd.pdb -E -j -R -P PIBOX_PUBLICIP -p 30000:30009
numprocs=1
autostart=true
autorestart=true
