# src/etc/supervisor/conf.d/torrent.conf
[program:torrent]
command=/usr/bin/rtorrent
numprocs=1
autostart=true
autorestart=true
